# deploy/systemd/caregist-pipeline-watchdog.service
[Unit]
Description=CareGist Pipeline Watchdog
Documentation=https://github.com/sapopo93/caregist
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=caregist
WorkingDirectory=/home/caregist/CareGist
ExecStart=/home/caregist/CareGist/.venv/bin/python /home/caregist/CareGist/tools/check_new_registration_pipeline.py --notify
StandardOutput=journal
StandardError=journal
SyslogIdentifier=caregist-pipeline-watchdog
EnvironmentFile=/home/caregist/CareGist/.env

[Install]
WantedBy=multi-user.target
