# deploy/systemd/caregist-incremental-update.service
[Unit]
Description=CareGist Incremental CQC Update
Documentation=https://github.com/sapopo93/caregist
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=caregist
WorkingDirectory=/home/caregist/CareGist
ExecStart=/home/caregist/CareGist/.venv/bin/python /home/caregist/CareGist/incremental_update.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=caregist-incremental-update
EnvironmentFile=/home/caregist/CareGist/.env

[Install]
WantedBy=multi-user.target
