# deploy/systemd/caregist-flush-email-queue.service
[Unit]
Description=CareGist Flush Email Queue
Documentation=https://github.com/sapopo93/caregist
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=caregist
WorkingDirectory=/home/caregist/CareGist
ExecStart=/home/caregist/CareGist/.venv/bin/python /home/caregist/CareGist/tools/flush_email_queue.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=caregist-flush-email-queue
EnvironmentFile=/home/caregist/CareGist/.env

[Install]
WantedBy=multi-user.target
