# deploy/systemd/caregist-incremental-update.timer
[Unit]
Description=CareGist Incremental CQC Update — every 30 minutes
Documentation=https://github.com/sapopo93/caregist
# Replaces: */30 * * * * in /etc/cron.d/caregist

[Timer]
OnCalendar=*:0/30
Persistent=true
AccuracySec=1min
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
